# dashboard/systemd/dashboard-agent.service
[Unit]
Description=Homelab Dashboard Stats Agent
After=network.target docker.service
Wants=network-online.target

[Service]
ExecStart=/usr/bin/python3 /home/clay/dashboard/agent.py
WorkingDirectory=/home/clay/dashboard
Restart=always
RestartSec=5
User=clay
# Pi-hole widget (the node running Pi-hole only): set the web/app password
# so the agent can query the admin API. Uncomment and fill in:
# Environment=PIHOLE_PASSWORD=your-pihole-password
# Environment=PIHOLE_BASE_URL=http://localhost

[Install]
WantedBy=multi-user.target
