# dashboard/systemd/dashboard-agent.service
[Unit]
Description=Homelab Dashboard Stats Agent
After=network.target docker.service
Wants=network-online.target

[Service]
ExecStart=/usr/bin/python3 /home/clay/dashboard/agent.py
WorkingDirectory=/home/clay/dashboard
Restart=always
RestartSec=5
User=clay
# Optional integrations — uncomment and fill in on the relevant node.
# Pi-hole widget (node running Pi-hole):
# Environment=PIHOLE_PASSWORD=your-pihole-password
# Environment=PIHOLE_BASE_URL=http://localhost
# Jellyfin now-playing (node running Jellyfin):
# Environment=JELLYFIN_API_KEY=your-api-key
# Environment=JELLYFIN_BASE_URL=http://localhost:8096
# Remndrs open-reminder count:
# Environment=REMNDRS_COUNT_URL=http://localhost:3000/api/reminders/open/count
# Meshtastic widget (node running meshtasticd / radio with TCP API):
# Environment=MESHTASTIC_HOST=localhost

[Install]
WantedBy=multi-user.target
